# systemd/video-studio-render-watcher.path
[Unit]
Description=Watch video-studio render trigger and run host-side render writer

[Path]
PathModified=/root/.openclaw/workspace/skills/video-studio/.video-render-trigger
Unit=video-studio-render-watcher.service

[Install]
WantedBy=multi-user.target
