# systemd/video-studio-narrate-watcher.path
[Unit]
Description=Watch video-studio narrate trigger and run host-side narrate writer

[Path]
PathModified=/root/.openclaw/workspace/skills/video-studio/.video-narrate-trigger
Unit=video-studio-narrate-watcher.service

[Install]
WantedBy=multi-user.target
