# systemd/video-studio-script-watcher.service
[Unit]
Description=Process pending video-studio script jobs (write narration drafts)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/.openclaw/workspace/skills/video-studio
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/root/.openclaw/workspace/skills/video-studio/scripts:/root/.openclaw/workspace/skills/voice-studio/scripts:/root/.openclaw/workspace/public-downloads
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/skills/video-studio/scripts/process_video_script_jobs.py
