# systemd/video-studio-script-watcher.path
[Unit]
Description=Watch video-studio script trigger and run host-side script writer

[Path]
PathModified=/root/.openclaw/workspace/skills/video-studio/.video-script-trigger
Unit=video-studio-script-watcher.service

[Install]
WantedBy=multi-user.target
